# rootfs/etc/services.d/web/run
#!/usr/bin/env bash
set -euo pipefail
cd /workspace
python3 -m http.server "${WEB_PORT:-8088}"
